# debian/postrm
#!/bin/sh -e

if [ "$1" = "purge" ]; then
	update-rc.d entente remove defaults >/dev/null
fi
